# packer/ansible/roles/vault_apiserver_network_proxy_server/templates/vault-apiserver-network-proxy-server.timer.j2
[Unit]
Description=Vault agent timer for apiserver-network-proxy

[Timer]
OnBootSec=0min
OnUnitActiveSec=20h
RandomizedDelaySec=1h

[Install]
WantedBy=timers.target
